# etc/supervisor/conf.d/nakadia.conf
[program:nakadia]
command=/usr/sbin/splunk
user=root

