# etc/supervisor/conf.d/nakadia.conf
[program:splunk]
command=/usr/sbin/splunk
user=root

